# templates/gmond-init.d.erb
#! /bin/sh
### BEGIN INIT INFO
# Provides:          ganglia-monitor.<%= @name %>
# Required-Start:    $network $named $remote_fs $syslog $all
# Required-Stop:     $network $named $remote_fs $syslog
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
### END INIT INFO 
PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin
DAEMON=/usr/sbin/gmond
NAME=<%= @name_real %>
DESC="Ganglia Monitor Daemon: <%= @name_real %>"
CONF="/etc/ganglia/<%= @name_real %>.conf"

test -x $DAEMON || exit 0

. /lib/lsb/init-functions

set -e

case "$1" in
  start)
        echo -n "Starting $DESC: "
        start-stop-daemon --start --quiet --pidfile /var/run/$NAME.pid \
                --exec $DAEMON -- --pid-file=/var/run/$NAME.pid --conf=$CONF
        echo "$NAME."
        ;;
  stop)
        echo -n "Stopping $DESC: "
        start-stop-daemon --stop  --quiet --oknodo --pidfile /var/run/$NAME.pid \
                2>&1 > /dev/null
        echo "$NAME."
        ;;
  status)
        status_of_proc -p "/var/run/${NAME}.pid" "${DAEMON}" "${NAME}"
        ;;
  reload)
  ;;
  restart|force-reload)
        $0 stop
        $0 start
        ;;
  *)
        N=/etc/init.d/$NAME
        echo "Usage: $N {start|stop|status|restart|force-reload}" >&2
        exit 1
        ;;
esac

exit 0
